# infra/systemd/orchestrator-stack.service
[Unit]
Description=FleetUpdate Stack (Server + UI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=orchestrator
Group=orchestrator
WorkingDirectory=/opt/orchestrator
Environment=PYTHONUNBUFFERED=1
Environment=SERVER_PSK=changeme
ExecStart=/bin/bash scripts/run-stack.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
