# infra/systemd/orchestrator-ui.service
[Unit]
Description=GitOps Orchestrator UI (Vite Dev/Preview)
After=network-online.target

[Service]
Type=simple
User=orchestrator
Group=orchestrator
# Change to your deployment path
WorkingDirectory=/opt/orchestrator
Environment=HOST=0.0.0.0
ExecStart=/bin/bash scripts/run-ui.sh
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
